# technitium-dns/rootfs/etc/s6-overlay/s6-rc.d/cert-convert-startup/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Generate PKCS #12 certificate on startup
# ==============================================================================

# Check if SSL certificate generation is enabled
if ! bashio::config.true 'generate_ssl_pfx'; then
    bashio::log.info "SSL certificate generation disabled - skipping"
    exit 0
fi

bashio::log.info "================================================="
bashio::log.info "  Add On - SSL Certificate Manager"
bashio::log.info "================================================="

source "/etc/s6-overlay/scripts/cert_utils.sh"

if ! check_pkcs12; then
    check_cert_paths
    if [ ! -f "$CERT_FILE" ] || [ ! -f "$KEY_FILE" ]; then
        generate_self_signed
    fi
    generate_pkcs12
fi
