# technitium-dns/rootfs/etc/s6-overlay/s6-rc.d/cert-convert-startup/up
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Generate PKCS #12 certificate from Let's Encrypt certificates on startup
# ==============================================================================

source "/etc/s6-overlay/scripts/cert_utils.sh"

bashio::log.info "SSL PFX generation check starting..."

if ! bashio::config.true 'generate_ssl_pfx'; then
    bashio::log.info "SSL PFX generation disabled"
    exit 0
fi

if ! command -v openssl >/dev/null 2>&1; then
    bashio::log.error "openssl binary not found!"
    exit 1
fi

check_and_generate || bashio::exit.nok "SSL PFX generation failed"

bashio::log.info "SSL PFX generation completed successfully"
