# technitium-dns/rootfs/etc/s6-overlay/s6-rc.d/cert-convert-watch/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Watch Let's Encrypt certificates for changes and regenerate PKCS #12 if needed
# ==============================================================================

source "/etc/s6-overlay/scripts/cert_utils.sh"

bashio::log.info "================================================="
bashio::log.info "  Starting certificate watch using inotify..."
bashio::log.info "================================================="

while true; do
    inotifywait -e close_write,move,create "${CERT_SOURCE_CHAIN}" "${CERT_SOURCE_KEY}" >/dev/null 2>&1
    bashio::log.info "Detected certificate change — checking and regenerating if needed..."
    if check_and_generate; then
        bashio::log.info "Restarting Technitium DNS Server to load new certificate..."
        s6-svc -r /run/s6-rc/supervise/dns-server
    fi
done
