# technitium-dns/rootfs/etc/s6-overlay/s6-rc.d/init-dns-server/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# shellcheck disable=SC1091
# ==============================================================================
# Set up environment variables for Technitium DNS Server
# ==============================================================================

# Enable strict error handling
set -o nounset -o errexit -o pipefail

bashio::log.info "ENV: Initializing Technitium DNS environment variables..."

# -----------------------------------------------------------------------------
# Module Dependencies
# -----------------------------------------------------------------------------
# Load utility modules
# shellcheck source=rootfs/etc/s6-overlay/scripts/utils/init_utils.sh
source /etc/s6-overlay/scripts/init_utils.sh

initialize_configuration
