# technitium-dns/rootfs/etc/s6-overlay/s6-rc.d/dns-server/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Starts Technitium DNS Server
# ==============================================================================

sleep 3

bashio::log.debug "================================================="
bashio::log.debug "  Starting Technitium DNS Server"
bashio::log.debug "================================================="

cd "/opt/technitium/dns" || bashio::exit.nok "Could not change directory to Technitium DNS Server"

command -v /usr/share/dotnet/dotnet >/dev/null 2>&1 || bashio::exit.nok "dotnet not found!"
bashio::log.debug "Dotnet version: $(/usr/share/dotnet/dotnet --version)"

exec /usr/share/dotnet/dotnet /opt/technitium/dns/DnsServerApp.dll /config
