# technitium-dns/rootfs/etc/s6-overlay/s6-rc.d/dns-server/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# shellcheck disable=SC1091
# ==============================================================================
# Starts Technitium DNS Server
# ==============================================================================

bashio::log.debug "================================================="
bashio::log.debug "  Starting Technitium DNS Server"
bashio::log.debug "================================================="

# -----------------------------------------------------------------------------
# Environment Validation
# -----------------------------------------------------------------------------
# Ensure we're in the correct working directory
cd "/opt/technitium/dns" || bashio::exit.nok "Could not change directory to Technitium DNS Server"

# Verify dotnet runtime is available
if ! command -v /usr/share/dotnet/dotnet >/dev/null 2>&1; then
    bashio::exit.nok "Required dependency 'dotnet' not found!"
fi

# -----------------------------------------------------------------------------
# Launch DNS Server
# -----------------------------------------------------------------------------
# Start the DNS server in background
/usr/share/dotnet/dotnet /opt/technitium/dns/DnsServerApp.dll /config
