# technitium-dns/rootfs/etc/s6-overlay/s6-rc.d/dns-server/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Starts Technitium DNS Server
# ==============================================================================

# Add a delay before starting for cleaner logs
if bashio::config.true 'generate_ssl_pfx'; then
    sleep 1
fi

bashio::log.info "================================================="
bashio::log.info "  Starting Technitium DNS Server"
bashio::log.info "================================================="

cd "/opt/technitium/dns" || bashio::exit.nok "Could not change directory to Technitium DNS Server"

command -v /usr/share/dotnet/dotnet >/dev/null 2>&1 || bashio::exit.nok "dotnet not found!"
bashio::log.info "Dotnet version: $(/usr/share/dotnet/dotnet --version)"

exec /usr/share/dotnet/dotnet /opt/technitium/dns/DnsServerApp.dll /config
