# technitium-dns/rootfs/etc/s6-overlay/s6-rc.d/dns-server/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Starts Technitium DNS Server
# ==============================================================================

# Add strict error handling
set -euo pipefail

# Configure error handling
trap 'handle_error $? $LINENO $BASH_LINENO "$BASH_COMMAND" $(printf "::%s" ${FUNCNAME[@]:-})' ERR

handle_error() {
    local exit_code=$1
    local line_no=$2
    local bash_lineno=$3
    local last_command=$4
    local func_trace=$5

    bashio::log.error "Error in bash at line ${bash_lineno}"
    bashio::log.error "Error in init-dns-server at line ${line_no}"
    bashio::log.error "Command '${last_command}' exited with status ${exit_code}"
    bashio::log.error "Function trace: ${func_trace}"

    exit "${exit_code}"
}

LOG_LEVEL=$(bashio::config 'log_level' 'info')
bashio::log.level "$LOG_LEVEL"

sleep 3

bashio::log.debug "================================================="
bashio::log.debug "  Starting Technitium DNS Server"
bashio::log.debug "================================================="

cd "/opt/technitium/dns" || bashio::exit.nok "Could not change directory to Technitium DNS Server"

command -v /usr/share/dotnet/dotnet >/dev/null 2>&1 || bashio::exit.nok "dotnet not found!"
bashio::log.debug "Dotnet version: $(/usr/share/dotnet/dotnet --version)"

exec /usr/share/dotnet/dotnet /opt/technitium/dns/DnsServerApp.dll /config
