# technitium-dns/rootfs/etc/s6-overlay/s6-rc.d/dns-server/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# shellcheck disable=SC1091
# ==============================================================================
# Home Assistant Add-on: Technitium DNS Server
# Manages the startup, readiness detection, and initialization of DNS server
# ==============================================================================

# Enable strict mode for better error handling and debugging
set -o nounset -o errexit -o pipefail

# -----------------------------------------------------------------------------
# Startup Banner
# -----------------------------------------------------------------------------
# Display a visual separator and startup message in debug logs
bashio::log.debug "dns-server: ================================================="
bashio::log.debug "dns-server:   Starting Technitium DNS Server"
bashio::log.debug "dns-server: ================================================="

# -----------------------------------------------------------------------------
# Environment Validation
# -----------------------------------------------------------------------------
# Ensure required dependencies and environment are ready before proceeding

# Change to the Technitium DNS application directory
cd "/opt/technitium/dns" || bashio::exit.nok "Could not change directory to Technitium DNS Server"

# Verify the .NET runtime is available as it's required by Technitium
if ! command -v /usr/share/dotnet/dotnet >/dev/null 2>&1; then
    bashio::exit.nok "Required dependency 'dotnet' not found!"
fi

# Disable .NET diagnostics
DOTNET_EnableDiagnostics=0

# -----------------------------------------------------------------------------
# Launch DNS Server
# -----------------------------------------------------------------------------
# Start the DNS server with the .NET runtime
exec /usr/share/dotnet/dotnet /opt/technitium/dns/DnsServerApp.dll /config
