# technitium-dns/rootfs/etc/s6-overlay/s6-rc.d/dns-server/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# shellcheck disable=SC1091
# ==============================================================================
# Home Assistant Add-on: Technitium DNS Server
# Manages the startup, readiness detection, and initialization of DNS server
# ==============================================================================

# Enable strict mode for better error handling and debugging
set -o nounset -o errexit -o pipefail

# -----------------------------------------------------------------------------
# Startup Banner
# -----------------------------------------------------------------------------
# Display a visual separator and startup message in debug logs
bashio::log.info "dns-server: ================================================="
bashio::log.info "dns-server:   Starting Technitium DNS Server"
bashio::log.info "dns-server: ================================================="

# -----------------------------------------------------------------------------
# Environment Validation
# -----------------------------------------------------------------------------
# Ensure required dependencies and environment are ready before proceeding

# Change to the Technitium DNS application directory
cd "/opt/technitium/dns" || bashio::exit.nok "Could not change directory to Technitium DNS Server"

# Verify the .NET runtime is available as it's required by Technitium
if ! command -v /usr/share/dotnet/dotnet >/dev/null 2>&1; then
    bashio::exit.nok "Required dependency 'dotnet' not found!"
fi

# Disable .NET diagnostics
DOTNET_EnableDiagnostics=0

# -----------------------------------------------------------------------------
# First Run Configuration
# -----------------------------------------------------------------------------
# Set up default DNS forwarders on first run
FIRST_RUN_MARKER="/data/addon-first-run-completed"
if [ ! -f "$FIRST_RUN_MARKER" ]; then
    bashio::log.debug "First run detected - configuring default DNS settings"

    # Set default DNS server name and domain
    export DNS_SERVER_DOMAIN="$(bashio::info.hostname)"

    # Set default forwarders (Cloudflare DoH)
    export DNS_SERVER_FORWARDERS="https://cloudflare-dns.com/dns-query (1.1.1.1),https://cloudflare-dns.com/dns-query (1.0.0.1)"
    export DNS_SERVER_FORWARDER_PROTOCOL="Https"

    # Create marker file to indicate first run setup is complete
    touch "$FIRST_RUN_MARKER"
    bashio::log.debug "Default DNS server domain set to: $DNS_SERVER_DOMAIN"
    bashio::log.debug "Default DNS forwarders configured: $DNS_SERVER_FORWARDERS (Protocol: $DNS_SERVER_FORWARDER_PROTOCOL)"
fi

# -----------------------------------------------------------------------------
# Launch DNS Server
# -----------------------------------------------------------------------------
# Start the DNS server with the .NET runtime
exec /usr/share/dotnet/dotnet /opt/technitium/dns/DnsServerApp.dll /config
