# technitium-dns/rootfs/etc/s6-overlay/s6-rc.d/ingress/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Technitium DNS Server
# Ingress reverse proxy: forwards the Supervisor-assigned ingress port to the
# Technitium web console on 127.0.0.1:5380.
#
# Required because the add-on uses host networking, for which Home Assistant
# assigns a dynamic ingress port (ingress_port: 0) that must be read at runtime.
# ==============================================================================

set -o nounset -o errexit -o pipefail

# Port assigned by the Supervisor for this add-on's ingress session
declare ingress_port
ingress_port=$(bashio::addon.ingress_port)

bashio::log.info "ingress: proxying port ${ingress_port} -> 127.0.0.1:5380"

# Render the server block with the assigned port
mkdir -p /etc/nginx/servers
cat > /etc/nginx/servers/ingress.conf <<EOF
server {
    listen ${ingress_port} default_server;
    server_name _;

    location / {
        proxy_pass http://127.0.0.1:5380;
        proxy_http_version 1.1;
        # Technitium sends X-Frame-Options: DENY and CSP frame-ancestors 'none',
        # which stop Home Assistant Ingress from embedding the console in its
        # iframe. Ingress serves the console same-origin with the HA frontend,
        # so relax framing to same-origin only (instead of removing the
        # protection entirely): this lets Ingress embed it while still blocking
        # cross-origin clickjacking, and keeps Technitium's script/style CSP.
        proxy_hide_header X-Frame-Options;
        proxy_hide_header Content-Security-Policy;
        add_header X-Frame-Options "SAMEORIGIN" always;
        add_header Content-Security-Policy "default-src 'self'; script-src 'self' 'unsafe-inline' 'unsafe-eval'; style-src 'self' 'unsafe-inline'; img-src 'self' data:; frame-ancestors 'self';" always;
        proxy_set_header Host \$host;
        proxy_set_header X-Real-IP \$remote_addr;
        proxy_set_header X-Forwarded-For \$proxy_add_x_forwarded_for;
        proxy_set_header X-Forwarded-Proto \$scheme;
        proxy_set_header Upgrade \$http_upgrade;
        proxy_set_header Connection \$connection_upgrade;
        proxy_read_timeout 3600s;
        proxy_buffering off;
    }
}
EOF

exec nginx -g "daemon off;" -c /etc/nginx/nginx.conf
