# technitium-dns/rootfs/etc/s6-overlay/s6-rc.d/cert-watch/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# shellcheck disable=SC1091
# shellcheck disable=SC2154
# ==============================================================================
# Watch certificates for changes and regenerate PKCS #12 if needed
# ==============================================================================

bashio::log.debug "================================================="
bashio::log.debug "  Starting Certificate Tool"
bashio::log.debug "================================================="

# -----------------------------------------------------------------------------
# Configuration Check
# -----------------------------------------------------------------------------
# Skip certificate monitoring if user prefers manual handling
if bashio::config.true 'manualCertHandling'; then
    bashio::log.debug "Manual certificate handling enabled - disabling certificate watch"
    exec sleep infinity
fi

# -----------------------------------------------------------------------------
# Source Dependencies
# -----------------------------------------------------------------------------
# Source utilities and perform initial setup
# shellcheck source=/etc/s6-overlay/scripts/cert_utils.sh
if ! source "/etc/s6-overlay/scripts/cert_utils.sh"; then
    bashio::exit.nok "Failed to source certificate utilities"
fi

# shellcheck source=/etc/s6-overlay/scripts/helper_utils.sh
if ! source "/etc/s6-overlay/scripts/helper_utils.sh"; then
    bashio::exit.nok "Failed to source helper utilities"
fi

# -----------------------------------------------------------------------------
# Helper Functions
# -----------------------------------------------------------------------------
# Return a list of unique directories that need to be monitored
get_watch_dirs() {
    # Get certificate file directories
    local cert_dir key_dir pkcs12_dir

    # Extract directories separately to avoid masking return values
    cert_dir=$(dirname "${CERT_FILE}" || true)
    key_dir=$(dirname "${KEY_FILE}" || true)
    pkcs12_dir=$(dirname "${PKCS12_FILE}" || true)

    # Combine and deduplicate the directories
    echo "${cert_dir}" "${key_dir}" "${pkcs12_dir}" | tr ' ' '\n' | sort -u
}

# Check for required inotify utility
check_inotify() {
    if ! command -v inotifywait >/dev/null 2>&1; then
        bashio::log.warning "inotifywait not found - certificate auto-updates may not work"
        return 1
    fi
    return 0
}

# -----------------------------------------------------------------------------
# Debounce Handling
# -----------------------------------------------------------------------------
# Prevent too frequent updates when multiple files change at once
readonly DEBOUNCE_SECONDS=10
handle_cert_changes() {
    local last_run=0
    local now

    # Get current timestamp
    now=$(date +%s)

    # Skip if we've updated recently
    if ((now - last_run < DEBOUNCE_SECONDS)); then
        bashio::log.debug "Skipping update (within debounce period)"
        return 0
    fi

    # Process certificate updates
    bashio::log.debug "Updating certificates..."
    handle_cert_update
    last_run=${now}
}

# -----------------------------------------------------------------------------
# Initial Setup
# -----------------------------------------------------------------------------
# Perform initial certificate setup
bashio::log.debug "Performing initial certificate setup..."
handle_cert_update

# Create all required directories to ensure they exist for monitoring
for dir in $(get_watch_dirs); do
    mkdir -p "${dir}"
done

# Wait for DNS server to fully initialize
sleep 10

# -----------------------------------------------------------------------------
# Certificate Monitoring Loop
# -----------------------------------------------------------------------------
# Start continuous monitoring for certificate changes
bashio::log.debug "Starting certificate watch service..."
while true; do
    # Check for proper inotify functionality
    if ! check_inotify; then
        bashio::log.warning "inotify check failed, falling back to periodic checks"
        sleep 300 # Check every 5 minutes as a fallback
        handle_cert_changes
        continue
    fi

    # Watch parent directories for any changes to certificate files
    IFS=$'\n' read -r -d '' -a watch_dirs < <(get_watch_dirs || true)
    if ! inotifywait -e close_write,move,create,delete "${watch_dirs[@]}" >/dev/null 2>&1; then
        bashio::log.debug "inotifywait failed, retrying in 10 seconds..."
        sleep 10
        continue
    fi

    # Handle detected changes
    bashio::log.debug "Detected certificate change, processing..."
    handle_cert_changes
done
